# pillars.service
# Pillars Character Generator - Systemd Service File
#
# Installation:
#   1. Copy to /etc/systemd/system/pillars.service
#   2. Update paths and user as needed
#   3. Run: sudo systemctl daemon-reload
#   4. Run: sudo systemctl enable pillars
#   5. Run: sudo systemctl start pillars
#
[Unit]
Description=Pillars Character Generator
After=network.target

[Service]
Type=simple
User=sam
Group=sam
WorkingDirectory=/home/sam/dev/pillars-character-gen/webapp
EnvironmentFile=/home/sam/dev/pillars-character-gen/.env
Environment="GUNICORN_WORKERS=3"
Environment="GUNICORN_BIND=0.0.0.0:9066"
Environment="PYTHONPATH=/home/sam/dev/pillars-character-gen"
ExecStart=/home/sam/dev/pillars-character-gen/.venv/bin/gunicorn \
    --workers $GUNICORN_WORKERS \
    --bind $GUNICORN_BIND \
    --access-logfile - \
    --error-logfile - \
    webapp.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
